# freebsd/rspamd.sh.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: rspamd
# REQUIRE: LOGIN
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable countd:
# rspamd (bool):          Set to "NO" by default.
#                          Set it to "YES" to enable rspamd.

. /etc/rc.subr

name="rspamd"
rcvar=`set_rcvar`
procname="@CMAKE_INSTALL_PREFIX@/bin/rspamd"
restart_precmd="rspamd_checkconfig"
reload_precmd="rspamd_checkconfig"
configtest_cmd="rspamd_checkconfig"

load_rc_config $name

: ${rspamd_enable="NO"}
: ${rspamd_pidfile="/var/run/rspamd/rspamd.pid"}
: ${rspamd_user="@RSPAMD_USER@"}
: ${rspamd_group="@RSPAMD_GROUP@"}

pidfile="$rspamd_pidfile"

stop_postcmd="rm -f $rspamd_pidfile"

rspamd_checkconfig()
{
	echo "Performing sanity check on rspamd configuration:"
	eval ${command} -t
}

extra_commands="reload configtest"
sig_reload="USR1"

command="$procname"
command_args="-u ${rspamd_user} -g ${rspamd_group} -c @CMAKE_INSTALL_PREFIX@/etc/rspamd.conf"

run_rc_command "$1"
